# ops/podman/frontend.container
[Unit]
Description=Kassa Frontend (nginx)
After=backend.service

[Container]
ContainerName=frontend
# In Produktion auf ein :<sha>-Tag pinnen (Rollback). NICHT auto-updaten.
Image=ghcr.io/maverick-2bit/kassa-frontend:latest
Network=kassa.network
PublishPort=80:80

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
