# ops/podman/kundendisplay.container
[Unit]
Description=Kassa Kundendisplay (nginx)
After=backend.service

[Container]
ContainerName=kundendisplay
# In Produktion auf ein :<sha>-Tag pinnen (Rollback). NICHT auto-updaten.
Image=ghcr.io/maverick-2bit/kassa-kundendisplay:latest
Network=kassa.network
PublishPort=8081:80

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
